# configs/rpi-4-files/etc/uci-defaults/90-profile-aliases
#!/bin/sh

# add some commonly used aliases
file=/etc/profile
grep -q "alias l" ${file} ||
    cat <<EOF >> ${file}
alias l='ls -lF'
alias la='ls -alF'
EOF
