# configs/rpi-4-files/etc/uci-defaults/90-patch-fwknopd
#!/bin/sh

# see https://github.com/openwrt/luci/issues/5073
file=/www/luci-static/resources/view/fwknopd.js
sed -i "s/s:option/s.option/" ${file}
